# quadlet/kaine-model-server.container
# SPDX-License-Identifier: LicenseRef-CAL-0.2
# Copyright (c) 2026 Kaine.One <kaine.one@tuta.com>
[Unit]
Description=KAINE model server (OpenAI-compatible language organ)
Requires=kaine-redis.service
After=kaine-redis.service

[Container]
ContainerName=kaine-model-server
Image=ghcr.io/ggml-org/llama.cpp:server-cuda
Network=kaine.network
PublishPort=127.0.0.1:11434:8080
Volume=kaine-models.volume:/models:ro
# Rootless NVIDIA via CDI — requires `nvidia-ctk cdi generate` (see README).
AddDevice=nvidia.com/gpu=0
SecurityLabelDisable=true
Exec=--model /models/organ/KAINE-Qwen3.5-4B-abliterated.Q4_K_M.gguf --alias kaineone/Qwen3.5-4B-abliterated-GGUF --host 0.0.0.0 --port 8080
HealthCmd=curl -fsS http://127.0.0.1:8080/v1/models
HealthInterval=15s
HealthRetries=10

[Service]
Restart=always

[Install]
WantedBy=default.target
